# modules/ertp/files/ertp-api.conf
# upstart script for ertp api play webapp

description     "(re)start ertp api webapp"

start on runlevel [2345]
stop on runlevel [06]

env LANG=en_GB.UTF-8
env USER=deploy
env GROUP=deploy
env HOME=/home/deploy/ertp-api-play-1.0-SNAPSHOT

pre-start script
        test -c /dev/null || { stop; exit 0; }
        chdir ${HOME}
        rm RUNNING_PID || true
end script

script
exec start-stop-daemon --start --chuid $USER:$GROUP --pidfile ${HOME}/RUNNING_PID --exec ${HOME}/start -- -Dhttp.port=9900 -Dlogger.resource=ertp-api-logger.xml
end script