# modules/ganglia/files/etc/init/gmetad.conf
description "gmetad upstart script"
author "Web ops team"

start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 5 20

expect daemon

exec /usr/sbin/gmetad
