# modules/golang/files/etc/profile.d/gopath.sh
#!/bin/sh

export GOPATH=/var/govuk/gopath
