# modules/logstash/templates/etc/init/logstash-server.conf.erb
description "Logstash server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 120

env HOME=/var/apps/logstash

exec start-stop-daemon --start \
                       --chdir "$HOME" \
                       --chuid logstash:logstash \
                       --pidfile "/var/run/logstash/logstash-server.pid" \
                       --make-pidfile \
                       --startas "/usr/bin/java" -- \
                         -jar logstash.jar \
                         agent \
                         -f /etc/logstash/logstash-server.conf \
                         --log /var/log/logstash/logstash-server.log \
                         -- web --backend elasticsearch://localhost:<%= transport_port %>/logstash-server
