# modules/vpnc/files/etc/init/vpnc.conf
description "vpnc vpn client daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 120

exec vpnc network --no-detach
