# modules/govuk_apt/files/etc/cron.daily/remove_unused_kernels
#!/bin/bash
set -e -o pipefail

PATH=/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ubuntu-unused-kernels | xargs -r apt-get -y purge
