# modules/govuk_ci/templates/jenkins-agent.conf.erb
start on runlevel [2345]
stop on runlevel [!2345]

script
        test -f /etc/default/locale && . /etc/default/locale || true
        export LANG
        exec start-stop-daemon --start -c jenkins --exec /usr/bin/java \
                --name jenkins-agent --  -jar <%= @swarm_client_dest %> \
                -username jenkins_agent -password <%= @agent_user_api_token %> \
                -name <%= @hostname %> -executors 2 \
                -master https://<%= @ci_master %>/ \
                -fsroot <%= @swarm_user_home %> -labels '<%= @labels %>'
