# modules/logstash/files/etc/init/elasticsearch.conf
description "Basic Upstart Script for Elasticsearch"
author "Web ops team"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /var/apps/elasticsearch-0.18.7/bin/elasticsearch -f
