# modules/govuk/templates/etc/cron.d/snapshot.erb
*/10 0-4,6-23 * * * <%= @snapshot_user %> /usr/local/bin/rummager-snapshot
