# modules/govuk/templates/etc/cron.d/snapshot.erb
*/10 * * * * <%= @snapshot_user %> /usr/local/bin/rummager-snapshot
