# modules/openconnect/templates/openconnect.conf.erb
description "openconnect vpn client daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 120

<% if %w{yes no}.include?(@dnsupdate) -%>
<% @dns_text="export DNS_UPDATE=" + @dnsupdate + "&&" %> 
<% end %>

exec <%= @dns_text %> echo <%= @password %> | openconnect <%= @gateway %>/github --no-dtls --script /etc/vpnc/vpnc-script --user <%= @user %> --passwd-on-stdin
