# modules/puppet/templates/puppetdb/upstart.conf.erb
description "PuppetDB centralized storage daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 120

script
  if [ -e /var/log/puppetdb/puppetdb-oom.hprof ]; then
    mv /var/log/puppetdb/puppetdb-oom.hprof /var/log/puppetdb/puppetdb-oom.hprof.prev
  fi

  exec start-stop-daemon --start \
                         --chuid puppetdb \
                         --pidfile /var/run/puppetdb.pid \
                         --make-pidfile \
                         --startas /usr/bin/java \
                         -- \
                         <%= @java_args %> \
                         -XX:OnOutOfMemoryError="kill -9 %p" \
                         -jar /usr/share/puppetdb/puppetdb.jar \
                         services -c /etc/puppetdb/conf.d \
                         >> /var/log/puppetdb/upstart.out.log \
                         2>> /var/log/puppetdb/upstart.err.log
end script
