# modules/govuk/files/etc/profile.d/govuk-guix.sh
#!/bin/sh

export PATH="$PATH:/var/govuk/govuk-guix/bin"
