# modules/guix/files/etc/init/guix-daemon.conf
description "Build daemon for GNU Guix"

start on runlevel [2345]

stop on runlevel [016]

exec /var/guix/profiles/per-user/root/guix-profile/bin/guix-daemon --build-users-group=guixbuild
