# modules/puppet/templates/etc/init/puppetmaster.conf.erb
# Puppet Master

description     "puppet master"

start on runlevel [2345]
stop on runlevel [!2345]

env PIDFILE=/var/run/puppetmaster/puppet-unicorn.pid
env CMD="unicornherder -u unicorn -p $PIDFILE -- -p <%= @unicorn_port %> -c /etc/puppet/unicorn.conf"

# If env HOME not set puppet is throwing the following error
#   Error: Could not intialize global default settings: couldn't find HOME
#   environment -- expanding ~
# See Puppet::Util::RunMode#which_dir which tries to expand ~ regardless of defaults
#  in config.ru
env HOME=/var/lib/puppet
<% if @puppet_sentry_dsn %>
env PUPPET_SENTRY_DSN=<%= @puppet_sentry_dsn %>
<% end %>

respawn
respawn limit 7 20

pre-start script
    # prepare environment
    mkdir -p /var/run/puppetmaster
    chown puppet:puppet /var/run/puppetmaster
end script

exec start-stop-daemon --start \
                       --chuid 'puppet:puppet' \
                       --chdir '/etc/puppet' \
                       --pidfile $PIDFILE \
                       --startas '/bin/sh' \
                       -- -c "exec $CMD"
