# modules/govuk/templates/apps/review_o_matic_explore/upstart.conf.erb
description "Proxy for Review-O-Matic Explore"

start on runlevel [2345]
stop on runlevel [!2345]

env PORT=3023
env UPSTREAM_AUTH=betademo:nottobeshared
env UPSTREAM_HOST="reviewomatic.<%= govuk_platform %>.alphagov.co.uk"

# Have to use this roundabout method of chuid'ing because lucid's upstart doesn't
# support the setuid and setgrp stanzas.
env USER=deploy
env GROUP=deploy

exec start-stop-daemon --start --chuid $USER:$USER --chdir /var/apps/review-o-matic-explore --exec /usr/bin/node server.js