# modules/govuk/templates/app_upstart.conf.erb
description "GOV.UK application: <%= name %>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env PORT=<%= port %>

script
  exec govuk_spinup "<%= name %>" "<%= app_type %>" \
    >>/var/log/<%= title %>/<%= title %>.log \
    2>>/var/log/<%= title %>/<%= title %>.error.log
end script
