# modules/govuk/templates/app_upstart.conf.erb
description "GOV.UK application: <%= name %>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env PORT=<%= port %>

exec govuk_spinup "<%= name %>" 2> /var/log/<%= title %>/<%= title.log %>
