# modules/govuk/templates/app_upstart.conf.erb
description "GOV.UK application: <%= title %>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

# If the app respawns more than 5 times in 20 seconds, it has deeper problems
# and should be killed off.
respawn limit 5 20

env PORT=<%= port %>

script
  exec govuk_spinup "<%= title %>" "<%= app_type %>" \
    >>/var/log/<%= title %>/upstart.out.log \
    2>>/var/log/<%= title %>/upstart.err.log
end script
