# modules/govuk/templates/procfile-worker.conf.erb
description "Procfile background worker manager for <%= @title %>"

start on runlevel [2345]
stop on runlevel [!2345]
<%- if !@enable_service -%>

manual
<%- end -%>

pre-start script
    for inst in $(seq 1 <%= @process_count%>)
    do
        start <%= @service_name%>_child INDEX=$inst
    done
end script

post-stop script
    for inst in $(initctl list | grep "^<%= @service_name%>_child " | awk '{print $2}' | tr -d ')' | tr -d '(')
    do
        stop <%= @service_name%>_child INDEX=$inst
    done

    # For backwards compatibility, also see if this job has a running process and kill it
    pids=$(initctl list | awk '/^<%= @service_name%> start\/running, process/ {print $4}')
    for pid in $pids
    do
      kill $pid
    done
end script
