# modules/govuk_cdnlogs/files/etc/cron.hourly/cdn_logs_rotate
#!/bin/sh
# cron job to rotate cdn logs.

/usr/sbin/logrotate /etc/logrotate.cdn_logs_hourly.conf
