# modules/mailhog/templates/mailhog.erb
#! /bin/sh
# /etc/init.d/mailhog
#
# MailHog init script.
#
# @author Jeff Geerling
# @author Ruben Arakelyan

### BEGIN INIT INFO
# Provides:          mailhog
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start MailHog at boot time.
# Description:       Enable MailHog.
### END INIT INFO

PID=/var/run/mailhog.pid
LOG=/var/log/mailhog.log
BIN=/usr/bin/mailhog
NOHUP_BIN=/usr/bin/nohup

# Carry out specific functions when asked to by the system
case "$1" in
  start)
    echo "Starting mailhog."
    if [ -f $PID ]; then
      echo "mailhog is already running.";
    else
      $NOHUP_BIN $BIN > $LOG 2>&1 &
      echo $! > $PID
    fi
    ;;
  stop)
    if [ -f $PID ]; then
      echo "Stopping mailhog.";
      kill -TERM $(cat $PID);
      rm -f $PID;
    else
      echo "mailhog is not running.";
    fi
    ;;
  restart)
    echo "Restarting mailhog."
    if [ -f $PID ]; then
      kill -TERM $(cat $PID);
      rm -f $PID;
    fi
    $NOHUP_BIN $BIN > $LOG 2>&1 &
    echo $! > $PID
    ;;
  status)
    if [ -f $PID ]; then
      echo "mailhog is running.";
    else
      echo "mailhog is not running.";
      exit 3
    fi
    ;;
  *)
    echo "Usage: /etc/init.d/mailhog {start|stop|status|restart}"
    exit 1
    ;;
esac

exit 0
