# modules/govuk/files/xvfb.conf
#!upstart

description "Virtual framebuffer xvfb"

start on local-filesystems
stop on shutdown

exec /usr/bin/Xvfb :99 -ac -screen 0 1024x768x16

respawn                # restart when job dies
respawn limit 5 60     # give up restart after 5 respawns in 60 seconds
